# scripts/com.taisei.translate-papers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.taisei.translate-papers</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ishimarutaisei/projects/ubbs/paper-translate/scripts/translate-papers-watch.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/ishimarutaisei/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/ishimarutaisei</string>
    <key>OLLAMA_HOST</key>
    <string>http://localhost:11434</string>
    <key>OLLAMA_MODEL</key>
    <string>minimax-m3:cloud</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/ishimarutaisei/projects/ubbs/paper-translate/.logs/launchd.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ishimarutaisei/projects/ubbs/paper-translate/.logs/launchd.stderr.log</string>
</dict>
</plist>
